# openbsd-port/files/rnsd.rc.sample
#!/bin/sh
#
# $OpenBSD: rnsd.rc.sample,v 1.4 2026/05/24 ... $ # Use CVS style version
# rc.d script for the Reticulum Network Stack daemon (rnsd)
# Packaged as akita-rnsd for the Akita Reticulum Master Node project

# --- Configuration ---
# These paths are substituted by the Makefile during installation
daemon="/usr/local/bin/rnsd"
rnstatus_cmd="/usr/local/bin/rnstatus"
pf_rule_generator="/usr/local/libexec/akita-rns/generate_pf_rules.py"
config_validator="/usr/local/libexec/akita-rns/validate_config.py"
# Python path is substituted by Makefile sed command below:
_python_cmd="/usr/local/bin/python${MODPY_VERSION}"

# Default configuration file path (can be overridden by akita_rnsd_config in rc.conf)
_akita_rnsd_config_default="/etc/reticulum/config"

# Default user to run as (created by the port, can be overridden by akita_rnsd_user)
_akita_rnsd_user_default="_reticulum"

# Default daemon flags (can be added to via akita_rnsd_flags)
_akita_rnsd_daemon_flags_default=""

# Syslog facility and priority for daemon logs
_akita_rnsd_syslog_facility="daemon.info"

# Name of the pf anchor used by this service
_pf_anchor_name="akita-rnsd"
_pf_anchor_file="/etc/pf.anchors/${_pf_anchor_name}"

# --- rc.subr Integration ---
. /etc/rc.d/rc.subr

# --- Service Variables ---
rc_bg=YES
rc_reload=YES

# Use service-specific variable names from rc.conf(8) if defined
_akita_rnsd_config=${akita_rnsd_config:-${_akita_rnsd_config_default}}
_akita_rnsd_user=${akita_rnsd_user:-${_akita_rnsd_user_default}}
_akita_rnsd_flags="${_akita_rnsd_daemon_flags_default} ${akita_rnsd_flags}"

# Standard OpenBSD rc.subr variables
daemon_flags="${_akita_rnsd_flags} -c \"${_akita_rnsd_config}\""
daemon_user="${_akita_rnsd_user}"
daemon_logger="${_akita_rnsd_syslog_facility}"
pexp="${daemon}.*${_akita_rnsd_config}"

# --- Service Functions ---

rc_pre() {
	# Check configuration file exists and is readable
	if [ ! -f "${_akita_rnsd_config}" ] || [ ! -r "${_akita_rnsd_config}" ]; then
		echo "Error: Configuration file ${_akita_rnsd_config} missing or unreadable."
		return 1
	fi
	
	# Verify daemon user
	if ! id -u "${daemon_user}" >/dev/null 2>&1; then
		echo "Error: User ${daemon_user} does not exist."
		return 1
	fi
	
	# Validate configuration using python helper if it exists
	if [ -x "${config_validator}" ] && ${_python_cmd} -c "import yaml" >/dev/null 2>&1; then
		if ! ${_python_cmd} "${config_validator}" "${_akita_rnsd_config}"; then
			echo "Error: Configuration validation failed."
			return 1
		fi
	fi
	
	# Generate PF rules
	generate_pf_rules || return 1
	return 0
}

rc_post() {
	echo "Flushing rules from pf anchor '${_pf_anchor_name}'."
	pfctl -a "${_pf_anchor_name}" -F rules >/dev/null 2>&1 || true
}

rc_status() {
	if rc_check; then
		local _pid=$(pgrep -U "${daemon_user}" -f "${pexp}" | xargs)
		echo "${name} is running as pid ${_pid}."
		echo "--- Reticulum Status ---"
		if [ -x "${rnstatus_cmd}" ]; then
			 if su "${daemon_user}" -c "${rnstatus_cmd}"; then
				 :
			 else
				 echo " Warning: Failed to execute '${rnstatus_cmd}'."
				 echo " Ensure 'sharing_scope: system' is set in ${_akita_rnsd_config}."
			 fi
		fi
		echo "------------------------"
	else
		echo "${name} is not running."
	fi
}

generate_pf_rules() {
	if [ ! -x "${pf_rule_generator}" ] || ! ${_python_cmd} -c "import yaml" >/dev/null 2>&1; then
		echo "Warning: PF rule generator or dependencies missing, skipping PF config."
		return 0
	fi

	echo "Generating PF rules for anchor '${_pf_anchor_name}'..."
	local _tmpdir="${TMPDIR:-/tmp}"
	local _tmp_pf_rules=$(mktemp "${_tmpdir}/${name}_pf.XXXXXX" 2>/dev/null)
	
	_cleanup_tmp() { rm -f "${_tmp_pf_rules}"; }
	
	if ${_python_cmd} "${pf_rule_generator}" "${_akita_rnsd_config}" > "${_tmp_pf_rules}" 2>&1; then
		if pfctl -a "${_pf_anchor_name}" -f "${_tmp_pf_rules}" 2>&1; then
			_cleanup_tmp
			return 0
		else
			echo "Error: Failed to load PF rules. Syntax error:"
			cat "${_tmp_pf_rules}" >&2
			_cleanup_tmp
			return 1
		fi
	else
		local _ret=$?
		echo "Error: Failed to generate PF rules."
		cat "${_tmp_pf_rules}" >&2
		_cleanup_tmp
		return 1
	fi
}

rc_cmd() {
	local cmd=$1
	shift

	case $cmd in
	checkconfig)
		rc_pre || return 1
		echo "Configuration checks passed."
		;;
	start)
		rc_start
		;;
	stop)
		rc_stop
		;;
	restart)
		rc_stop
		rc_start
		;;
	reload)
		if ! rc_check; then
			echo "Error: ${name} is not running."
			return 1
		fi
		echo "Reloading ${name} configuration and pf rules..."
		rc_pre || return 1
		generate_pf_rules || return 1
		if pkill -HUP -U "${daemon_user}" -f "${pexp}" 2>/dev/null; then
			echo "SIGHUP sent successfully."
		else
			echo "Warning: Failed to send SIGHUP."
		fi
		;;
	status)
		rc_status
		;;
	*)
		echo "Usage: $0 {start|stop|reload|restart|status|checkconfig}"
		return 1
		;;
	esac
}

# --- Main Execution ---
rc_cmd $1
